# migration/migration-publish.cron
# /etc/cron.d/migration-publish.cron: crontab entry to migrate published Mediapackages per series



*/2 * * * *   root	/usr/bin/python /mnt/opencast3/migration/migrate_series.py -i @/mnt/opencast3/migration/serieslist.txt publish >> /mnt/opencast3/migration/cron-publish.log 2>&1 
