# migration/migration.cron
# /etc/cron.d/migration.cron: crontab entry to migrate mediapackages, one at a time


ROOT_DIR=/mnt/opencast/oc-migration
LOGFILE=$ROOT_DIR/cron.log
EXEC=$ROOT_DIR/migrate.sh

SHELL=/bin/bash


* * * * *   root	$SHELL $EXEC &>> $LOGFILE
