# migration/oc-config/migration.cron
# /etc/cron.d/migration-archive.cron: crontab entry to migrate archived Mediapackages per series



*/2 * * * *   root	/usr/bin/python /mnt/opencast3/migration/migration_series.py -i @/mnt/opencast3/migration/serieslist.txt /mnt/opencast3/migration/migrated_packages 2> /dev/null
